# C/crontab
* * * * * bash ~/Bash/C/bkp ~/Bash/C/a ~/Bash/C/b
